# scripts/sauron.service
[Unit]
After=multi-user.target

[Service]
ExecStart=/home/pi/sauron/run.sh
WorkingDirectory=/home/pi/sauron

[Install]
WantedBy=default.target
