# scripts/sauron-photo.service
[Unit]
After=multi-user.target

[Service]
ExecStart=/home/pi/sauron/takePhoto.sh
WorkingDirectory=/home/pi/sauron

[Install]
WantedBy=default.target
